# files/usr/lib/systemd/system/restic-manager.timer
[Unit]
Description=Daily Backup timer

[Timer]
OnCalendar=*-*-* 03:00:00
Unit=restic-manager.service

[Install]
WantedBy=multi-user.target

